# Projects/YourScreenServer/supervisord.conf
[program:ngrok]
command=ngrok http 80 --log stdout
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
